# backend/templates/etc/systemd/system/openvpn-net.service
[Unit]
Description=OpenVPN NAT iptables

[Service]
ExecStart=/usr/local/bin/openvpn-net.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
