# vm-setup/cloud-route.service
[Unit]
Description=Cloud Route Backend Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/cloud-route/backend
ExecStart=/bin/bash -c ". ./start.sh"
Restart=on-failure
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target